# deploy/systemd/shortsflow-automation.timer.in
[Unit]
Description=Run ShortsFlow automation daily

[Timer]
OnCalendar=*-*-* 02:00:00 America/Sao_Paulo
Persistent=true
Unit=shortsflow-automation.service

[Install]
WantedBy=timers.target
