# deploy/systemd/shortsflow-analytics-sync.service.in
[Unit]
Description=ShortsFlow daily YouTube Analytics performance collection
After=network-online.target tailscaled.service shortsflow-hub.service
Wants=network-online.target tailscaled.service

[Service]
Type=oneshot
WorkingDirectory=__REPO_ROOT__
Environment=SHORTSFLOW_REPO_ROOT=__REPO_ROOT__
ExecStart=__REPO_ROOT__/.venv/bin/python -m app.cli analytics-sync-run
