# deploy/systemd/yts-render-automation.timer.in
[Unit]
Description=Run YTS Render automation daily

[Timer]
OnCalendar=*-*-* 02:00:00 America/Sao_Paulo
Persistent=true
Unit=yts-render-automation.service

[Install]
WantedBy=timers.target
